# deploy/pm1-beads/bd-board.service
[Unit]
Description=Beads Project Board (read-only, tailnet only)
After=network-online.target tailscaled.service dolt-sql-server.service
# Wants (not Requires): a running board must survive a dolt-sql-server
# restart/outage and keep serving last-good with the stale banner (spec C7).
# Requires= would propagate dolt's deactivation and kill the board, which
# directly defeats C7. Cold-start still gates on Dolt via ExecStartPre (C5).
Wants=network-online.target dolt-sql-server.service
# Cap restart hammering on the shared box if the start gate keeps failing.
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
User=admin
# pm1 must be a server-mode bd client: a one-time
#   mkdir ~/beads-workspace && cd ~/beads-workspace && \
#   bd init --server --server-host 100.85.126.95 --server-port 3307 \
#     --server-user beads --database beads_workspace
# creates this dir + its .beads metadata (connects to the EXISTING shared
# DB; non-destructive). `bd board` resolves the remote DB via this .beads.
WorkingDirectory=/home/admin/beads-workspace
# systemd EnvironmentFile needs plain KEY=value (NO shell `export `). This is
# beads-client.env with `export ` stripped: sed 's/^export //' ... mode 600.
EnvironmentFile=/home/admin/beads-board.env
# C5: a (re)start does not proceed until the tailnet IP is up AND Dolt
# answers a SQL ping. (A board that is already running is unaffected by
# this gate and keeps serving last-good through Dolt blips — see C7.)
ExecStartPre=/bin/sh -c 'ip addr show | grep -qF 100.85.126.95 || (echo "tailnet IP missing" >&2; exit 1)'
ExecStartPre=/bin/sh -c '/usr/local/bin/dolt --host 100.85.126.95 --port 3307 --user "$BEADS_DOLT_SERVER_USER" --password "$BEADS_DOLT_PASSWORD" --no-tls sql -q "SELECT 1" >/dev/null 2>&1 || (echo "dolt ping failed" >&2; exit 1)'
# Base ticketing workspace is explicit; every project workspace is auto-discovered
# by the glob, expanded live on each fetch — register.sh just creates the dir and
# it appears on the board within the cache TTL, no edit/restart here. bd board runs
# read-only against each (no ReadWritePaths needed for globbed workspaces).
ExecStart=/usr/local/bin/bd serve-board --addr 100.85.126.95:8099 --refresh 30 --cache-ttl 20 --exec-timeout 10 \
  --workspace /home/admin/beads-workspace \
  --workspace-glob '/home/admin/beads-*-workspace'
Restart=always
RestartSec=5
# A single OOM stops the unit instead of entering a SIGKILL/restart loop
# that would hammer the shared host (co-exists with the live trading stack).
OOMPolicy=stop
# C6: bound footprint on the shared box. 384M (not 192M) leaves headroom
# for the Go server PLUS a forked `bd board` doing a Dolt query and the
# buffered JSON, so light concurrency does not OOM-loop.
MemoryMax=384M
CPUQuota=25%
TasksMax=64
Nice=10
IOSchedulingClass=idle
LimitNOFILE=4096
# Defense in depth: filesystem read-only EXCEPT the server-mode client
# workspace, where bd maintains runtime .beads state (locks, proxy dir,
# local_version, jsonl). Without this carve-out the strict sandbox blocks
# bd board's normal client writes.
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/admin/beads-workspace /home/admin/beads-creator-kb-factory-workspace /home/admin/beads-growthcapital-content-workspace /home/admin/beads-KreatorFlow-workspace /home/admin/beads-antigravity-dashboard-workspace
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
